# deploy/sfiDashMonitoring.service
[Unit]
Description=sfiDashMonitoring - frontend + backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/youruser/Projets/sfiDashMonitoring
# Ensure start.sh is executable
ExecStart=/home/youruser/Projets/sfiDashMonitoring/start.sh sfimonitoring.com
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
# Optional: set these if needed
# Environment=FRONTEND_DOMAIN=sfimonitoring.com
# Environment=FRONTEND_PORT=5173
# Environment=BACKEND_PORT=3001

[Install]
WantedBy=multi-user.target
