# deployed/services/sfiDashMonitoring-frontend.service
[Unit]
Description=SFI Dashboard Monitoring Frontend (Nginx)
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/nginx
ExecReload=/usr/sbin/nginx -s reload
ExecStop=/usr/sbin/nginx -s quit
PIDFile=/var/run/nginx.pid
Restart=on-failure
RestartSec=10

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/nginx /var/cache/nginx /var/run/nginx.pid

[Install]
WantedBy=multi-user.target
